# organize_security_files.timer
[Unit]
Description=Run the organize_security_files.service every day

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
